# repack/overlay/etc/systemd/system/containerd.service.d/mnt-data.conf
[Unit]
# containerd state lives on /mnt/data — start after it mounts, stop before it unmounts.
RequiresMountsFor=/mnt/data
